# srcs/phpmyadmin/srcs/supervisor.conf
[supervisord]
nodaemon=true

[program:phpmyadmin]
command = sh -c "php -S 0.0.0.0:5000 -t  /www/phpmyadmin  & sleep infinite && kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "telegraf && kill -s SIGTERM $(cat supervisord.pid)"